# etc/profile.d/jade-profile.sh
#!/bin/sh

# set default applications for the dashboard

export JADE_TERMINAL=deepin-terminal
export JADE_FILE_MANAGER=nautilus
export JADE_SOFTWARE_INSTALLER=pamac-manager
export JADE_SETTINGS=manjaro-settings-manager
export JADE_BROWSER=chromium
export JADE_DISTRIBUTOR_WELCOME_SCREEN=manjaro-welcome
