# etc/profile.d/jade-profile.sh
#!/bin/sh

# set default applications for the dashboard

export JTERMINAL=deepin-terminal
export JFILEMANAGER=nautilus
export JSOFTWAREINSTALLER=pamac-manager
export JSETTINGS=manjaro-settings-manager
export JBROWSER=chromium
